# webhooker/webhooker.conf
# /etc/init/webhooker.conf

start on started network
stop on stopping network

respawn
expect daemon

exec sudo -u webhooker "/usr/sbin/webhooker"
